# dist/rocky_9_3/installation_kits/solr/solr.service
[Unit]
Description=Apache SOLR solr-8.11.2

[Service]
Type=forking
User=solr
Environment=SOLR_INCLUDE=/usr/lib/x14/solr/solr-8.11.2/bin/solr.in.sh
ExecStart=/usr/lib/x14/solr/solr-8.11.2/bin/solr start
ExecStop=/usr/lib/x14/solr/solr-8.11.2/bin/solr stop
Restart=on-failure
LimitNOFILE=65000
LimitNPROC=65000
TimeoutSec=180s

[Install]
WantedBy=multi-user.target
