# dist/ansible/roles/spark/templates/spark-slave.service.j2
[Unit]
Description=Apache Spark {{ spark_version }} Slave scala 12
Wants=network-online.target
After=network-online.target

[Service]
User={{ spark_user }}
Group={{ org_group }}
Type=forking
SuccessExitStatus=143

WorkingDirectory={{ spark_home }}
ExecStart="{{ spark_home }}/sbin/start-slave.sh {{ spark_master_url }}"
ExecStop="{{ spark_home }}/sbin/stop-slave.sh {{ spark_master_url }}"

[Install]
WantedBy=multi-user.target
