# dist/ansible/roles/spark/templates/spark-master.service.j2
[Unit]
Description=Apache Spark {{ spark_version }} Slave scala ?
Wants=network-online.target
After=network-online.target

[Service]
User={{ spark_user }}
Group={{ spark_group }}
Type=forking
SuccessExitStatus=143

WorkingDirectory={{ spark_home }}
EnvironmentFile={{ spark_home }}/etc/env
ExecStartPost=/bin/bash -c "echo $MAINPID > {{ spark_home }}/etc/spark-slave.pid"
ExecStart="{{ spark_home }}/sbin/start-master.sh"
ExecStop="{{ spark_home }}/sbin/stop-master.sh"


[Install]
WantedBy=multi-user.target



