# dist/fedora_40/installation_kits/kyuubi/kyuubi.service
[Unit]
Description=Kyuubi 1.9.2
Wants=network-online.target
After=network-online.target

[Service]
User=kyuubi
Group=x14
Type=forking
SuccessExitStatus=0

WorkingDirectory=/usr/lib/x14/kyuubi/apache-kyuubi-1.9.2-bin
ExecStartPost=/bin/bash -c "echo $MAINPID > /var/lib/x14/kyuubi/kyuubi.pid"
ExecStart=/usr/lib/x14/kyuubi/apache-kyuubi-1.9.2-bin/bin/kyuubi start
ExecStop=/usr/lib/x14/kyuubi/apache-kyuubi-1.9.2-bin/bin/kyuubi stop

[Install]
WantedBy=multi-user.target
