# dist/fedora_40/installation_kits/schema-registry/schema-registry.service
[Unit]
Description=Schema registry 7.4.0
Wants=network-online.target
After=network-online.target

[Service]
User=kafka
Group=x14
Type=simple
SuccessExitStatus=0

WorkingDirectory=/usr/lib/x14/kafka/schema-registry
Environment="LOG_DIR=/var/lib/x14/kafka/schema-registry/log"
ExecStartPost=/bin/bash -c "echo $MAINPID > /var/lib/x14/kafka/schema-registry.pid"
ExecStart=java -jar kafka-schema-registry-package-7.4.0-standalone.jar schema-registry.json

[Install]
WantedBy=multi-user.target
