# dist/ubuntu_22_04/installation_kits/nessie_iceberg/nessie.service
[Unit]
Description=Project Nessie 0.54.0
Wants=network-online.target
After=network-online.target

[Service]
User=nessie
Group=x14
Type=simple
Restart=always
RestartSec=5

WorkingDirectory=/usr/lib/x14/nessie/
#EnvironmentFile=/usr/lib/x14/nessie/etc/env
ExecStart=/usr/lib/x14/nessie/nessie-quarkus-0.54.0-runner -Dnessie.version.store.type=ROCKS -Dnessie.version.store.rocks.db-path=/tmp/rocks-nessie -Dquarkus.http.host=localhost

[Install]
WantedBy=multi-user.target
