# dist/bash/rocky_9_1/installation_kits/nessie_delta/nessie.service
[Unit]
Description=Project Nessie 0.60.1
Wants=network-online.target
After=network-online.target

[Service]
User=nessie
Group=x14
Type=simple
Restart=always
RestartSec=5

WorkingDirectory=/usr/lib/x14/nessie/
#EnvironmentFile=/usr/lib/x14/nessie/etc/env
ExecStart=/usr/lib/x14/nessie/nessie-quarkus-0.60.1-runner -Dnessie.version.store.type=ROCKS -Dnessie.version.store.rocks.db-path=/var/lib/x14/nessie/rocks-nessie -Dquarkus.http.host=10.1.1.93

[Install]
WantedBy=multi-user.target
