# dist/rocky_9_1/installation_kits/superset/superset.service
[Unit]
Description=Superset daemon
After=syslog.target network.target

[Service]
User=superset
Type=oneshot
ExecStart=/usr/lib/x14/superset/start.sh
#ExecStop=/usr/lib/x14/zeppelin/zeppelin-0.10.1-bin-all/bin/zeppelin-daemon.sh stop
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
