# dist/ansible/roles/kyuubi/templates/kyuubi.service.j2
[Unit]
Description=Kyuubi {{ kyuubi_version }}
Wants=network-online.target
After=network-online.target

[Service]
User={{ kyuubi_user }}
Group={{ org_group }}
Type=forking
SuccessExitStatus=0

WorkingDirectory={{ kyuubi_home }}
ExecStart={{ kyuubi_home }}/bin/kyuubi start
ExecStop={{ kyuubi_home }}/bin/kyuubi stop

[Install]
WantedBy=multi-user.target
