# dist/ubuntu_22_04/installation_kits/nginx/install.sh
export NEEDRESTART_MODE=a

apt-get -y install nginx
cp *.cfg /etc/nginx/sites-enabled
systemctl enable nginx
systemctl start nginx

