# dist/rocky_8_7/installation_kits/nginx/install.sh

yum -y install nginx
systemctl enable nginx
#systemctl start nginx

